# aux_files/project-spec/meta-user/recipes-apps/r5ctrlr-scpi-serv/files/r5ctrlr-scpi-serv.service
[Unit]
Description=r5ctrlr-scpi-serv
  
[Service]
ExecStart=/usr/bin/r5ctrlr-scpi-serv
StandardOutput=journal+console
  
[Install]
WantedBy=multi-user.target
